# etc/supervisord-celery.conf
[supervisord]
nodaemon=true

[program:app-celery]
command=celery -A danbi worker -l info
