# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:app-nginx]
command=nginx

[program:app-uwsgi]
command=uwsgi --ini /home/service/etc/uwsgi.ini

[program:app-celery]
command=celery -A danbi worker -l info
